# cd-rip/cd-rip-notify.timer
[Unit]
Description=CD rip completion notifier — run every minute

[Timer]
OnCalendar=minutely
Persistent=false

[Install]
WantedBy=timers.target
