# openclaw/openclaw-gateway-restart.timer
[Unit]
Description=Restart openclaw-gateway daily at 4:00 UTC

[Timer]
OnCalendar=*-*-* 19:00:00
Persistent=true

[Install]
WantedBy=timers.target
